# supervisord.conf
[supervisord]
nodaemon=true

[program:uvicorn]
command=/usr/local/bin/uvicorn CRUD.asgi:application --port 8448 --reload
directory=/app
autorestart=true
startretries=3
user=root